# hookshot.service
# /etc/systemd/system/hookshot.service
[Unit]
Description=GitHub webhook listener for deploy
After=network.target

[Service]
ExecStart=/usr/local/bin/webhook -hooks /home/worker/hooks.json -port 9000
WorkingDirectory=/home/worker
Restart=always
User=worker

[Install]
WantedBy=multi-user.target
